# containerized/monitor/supervisord.conf
[supervisord]
nodaemon=true

[program:cAdvisor]
command=/bin/bash -c "xxxxx"
